# deploy/systemd/signalops-status-cache.service
[Unit]
Description=SignalOps status cache refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/signalops-status
ExecStart=/usr/bin/php /var/www/signalops-status/scripts/refresh-cache.php
