# deploy/systemd/signalops-status-cache.timer
[Unit]
Description=Warm the SignalOps status cache

[Timer]
OnBootSec=20s
OnUnitActiveSec=90s
AccuracySec=5s
RandomizedDelaySec=5s
Unit=signalops-status-cache.service

[Install]
WantedBy=timers.target
